# deploy/telegram-video-relay.service
[Unit]
Description=Telegram Video Relay Bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/opt/tg-video-relay-bot
EnvironmentFile=/opt/tg-video-relay-bot/.env
ExecStart=/opt/tg-video-relay-bot/.venv/bin/python -m tg_video_relay_bot
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
